# target/etc/systemd/system/find-and-mount.service
[Unit]
Description=Find local harddrive and mount it on /home/contest

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/find_and_mount start
ExecStop=/usr/local/bin/find_and_mount stop

[Install]
WantedBy=multi-user.target
